# etc/grooved.service

[Unit]
Description=Groovy music player daemon